# config/autostart/nitrogen.desktop
[Desktop Entry]
Type=Application
Exec=nitrogen --restore
Hidden=false
Name=Nitrogen
