# config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Exec=picom
Hidden=false
Name=Picom
