# config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Exec=picom -b
Hidden=false
Name=Picom
