# config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=polybar -r
Hidden=false
Name=Polybar
