# config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=polybar
Hidden=false
Name=Polybar
