# config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=~/.config/polybar/launch-polybar.sh
Hidden=false
Name=Polybar
